# systemd/auditor-client.service
[Unit]
Description=Auditor client

[Service]
Type=simple
ExecStart=/usr/sbin/auditor -C

[Install]
WantedBy=multi-user.target
