# systemd/auditor-server.service
[Unit]
Description=Auditor server

[Service]
Type=simple
ExecStart=/usr/sbin/auditor -S

[Install]
WantedBy=multi-user.target
